# etc/services.d/socat/run
#!/usr/bin/with-contenv sh

if [ -z "${BASESTATIONHOST}" ]; then
  sleep 3600  
else
  socat -d -d TCP-LISTEN:30003,fork TCP:"${BASESTATIONHOST}":"${BASESTATIONPORT}"
fi

