# etc/services.d/php-fpm/run
#!/usr/bin/with-contenv sh
/usr/sbin/php-fpm7.3 --nodaemonize 2>&1 | awk '{print "[php-fpm] " $0}'

