# etc/services.d/nginx/run
#!/usr/bin/with-contenv sh

set -eo pipefail
/usr/sbin/nginx 2>&1 | awk '{print "[nginx] " $0}'

