# etc/services.d/update_db/run
#!/usr/bin/with-contenv sh
set -eo pipefail
sleep 1h
s6-setuidgid flightairmap php /var/www/flightairmap/htdocs/scripts/update_db.php 2>&1 | awk '{print "[update_db] " $0}'
